# systemd/user/batterynotify.service
[Unit]
Description=Battery notify

[Service]
ExecStart=/bin/bash .local/bin/batterynotify.sh

